# roles/_lib/templates/service_configuration.conf
description '{{ service_file }}'

start on (local-filesystems and net-device-up)
stop on shutdown

respawn
respawn limit 5 60

script
  cd {{ chdir }}
  exec {{ command }}
end script
